# src/com/ubuntu/content/service/content-hub.service
[Unit]
Description=Content Hub Service

[Service]
Type=dbus
BusName=com.ubuntu.content.dbus.Service
ExecStart=/usr/bin/content-hub-service
